# files/etc/cron.d/acme
00 15 * * * ${acme_user} lockf -t 0 /tmp/acme-cron.lock acme.sh --cron --home ${acme_home} --log /dev/null --syslog 6 > /dev/null
