# files/etc/cron.d/davical-update-groups
MAILTO=root
0,15,30,45 * * * * ${davical_user} ${davical_repo_dir}/scripts/cron-sync-ldap.php $fqdn
