# files/etc/cron.d/invidious-update
MAILTO=root
0 3 * * * root /usr/local/libexec/invidious-update -q ${invidious_user} ${invidious_repo_dir}
