# files/etc/profile.d/locale.sh
export LANG=${locale}
export CHARSET=UTF-8
